# systemd/oom_guard.service
[Unit]
Description=OOM Guard - Memory Monitor Daemon
Documentation=https://github.com/marvinbraga/oom_guard
After=network.target

[Service]
Type=simple
# Configuration via CLI
# Adjust thresholds as needed: -m MEM_WARN,MEM_KILL -s SWAP_WARN,SWAP_KILL
ExecStart=/usr/local/bin/oom_guard -m 10,5 -s 10,5 -n -r 3600 --set-priority=-20

# Alternative: Use environment file for configuration
# EnvironmentFile=/etc/default/oom_guard
# ExecStart=/usr/local/bin/oom_guard

Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security hardening
CapabilityBoundingSet=CAP_KILL CAP_DAC_OVERRIDE CAP_SYS_NICE CAP_SYS_PTRACE
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/proc

# Resource limits
CPUQuota=10%
MemoryMax=50M

[Install]
WantedBy=multi-user.target
